# packaging/linux/deb/noetica.service
[Unit]
Description=Noetica agent-machine — AI dialogue backend (:8080)
Documentation=https://github.com/SocioProphet/Noetica
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/lib/noetica/noetica-launch
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NOETICA_AM_PORT=8080
Environment=OLLAMA_HOST=http://127.0.0.1:11434

[Install]
WantedBy=multi-user.target
